# rootfs/etc/services.d/resol/run
#!/usr/bin/with-contenv bash

exec python3 /root/run.py